# deploy/koo-mcp-controller.service
[Unit]
Description=koo-mcp controller
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
# 전용 서비스 유저. 핵심: docker 그룹에 넣지 않는다 → 컨트롤러가 털려도
# raw /var/run/docker.sock에 직접 못 붙고, 반드시 socket-proxy(2375)만 경유.
# /home/ubuntu(750) 진입을 위해 ubuntu 그룹만 보조로 부여.
User=koomcp
SupplementaryGroups=ubuntu
WorkingDirectory=/home/ubuntu/koomcp/controller
# 환경 고유값(AUTHKIT_DOMAIN, BASE_URL 등)은 커밋되지 않는 파일에서 읽는다.
# 배포 시 deploy/koomcp.env.example → /etc/koomcp.env 로 복사해 채워둘 것.
EnvironmentFile=/etc/koomcp.env
# Docker는 socket-proxy를 통해서만 호출 (raw docker.sock 미사용)
Environment=DOCKER_HOST=tcp://127.0.0.1:2375
Environment=RUNNER_IMAGE=web-surf-runner:latest
Environment=RUNNER_RUNTIME=runsc
Environment=RUNNER_NETWORK=surf-egress
ExecStart=/home/ubuntu/koomcp/controller/.venv/bin/python controller.py
Restart=always
RestartSec=3

# --- 프로세스 하드닝 ---
NoNewPrivileges=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
